# hash-balancer.service
[Unit]
Description=Hash Load Balancer Service
After=network.target

[Service]
User=root
WorkingDirectory=/root
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/root
Environment=MORPH_API_KEY=morph_fdBw4OOQ9NwU6REhYRtmnv
ExecStart=/venv/bin/uvicorn load_balancer:app --host 0.0.0.0 --port 8000 --log-level debug
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
